# slowmovie.timer
[Unit]
Description=Slow Movie Player Timer

[Timer]
OnUnitActiveSec=2m
AccuracySec=10s

[Install]
WantedBy=default.target
